# profile.d/00-kmoore.sh

############################################################
# Update the following to personalize your bash_profile
# this file will be automatically be sourced via /Users/kmoore/.bash_profile

## GPG Email
#The email address associated with your GPG key
export CORP_GPGKEY_ADDRESS="kevin@aurora.tech"

export INITIALS="KGM"

#Folder path to where you will be checking out git projects
export CORP_HOME="${HOME}/src/aurora"


#Folder path for shared onboarding docs
export CORP_ONBOARD="/Users/kmoore/Aurora/Onboarding"

#Path to your ipsec secrets file
export IPSEC_SECRETS_FILE=/usr/local/etc/ipsec.secrets

export CORP_USERNAME="kevin"

export KEY_SUFFIX="aurora.tech"
export GIT_ORG="[GITORG]"

# Set architecture flags for Ruby RVM to play nice
export ARCHFLAGS="-arch x86_64"

##########################################################

## Feel free to add your own shell customizations here

