# profile.d/05-mynav.sh
alias mtl="cd ${CORP_HOME}/metal"
alias ops='cd ${CORP_HOME}/it-ops/'
alias ppt='cd ${CORP_HOME}/puppet/'
alias av='cd ${CORP_HOME}/av/'
