# profile.d/70-phabricator.sh
#!/bin/bash

[[ -e "${CORP_HOME}/arcanist/arcanist/resources/shell/bash-completion" ]] && source "${CORP_HOME}/arcanist/arcanist/resources/shell/bash-completion"
