# profile.d/json-lint.sh
alias bourne="python -mjson.tool"
alias gr-json-lint="python -mjson.tool"
