# profile.d/zz-aws-env.lnk
/Users/kmoore/src/grnds/engineering/bash/aws-env.sh