# profile.d/zz-gr-aws-cfn-helpers.lnk
/Users/kmoore/src/grnds/engineering/bash/gr-aws-cfn-helpers.sh