# profile.d/zz-gr-share-file.lnk
/Users/kmoore/src/grnds/engineering/bin/gr-share-file.sh